# supervisord/phpfpm.conf
[program:php-fpm]
command=/usr/sbin/php-fpm -F
autorestart=true
